# systemd/rapido-init.target
[Unit]
Description=Rapido Init Default Target
AllowIsolate=yes
